# docker-compose/web/crontab
* * * * * cd /var/www && php artisan schedule:run >> /dev/output 2>&1
